# deploy/review-helper.service
[Unit]
Description=review-helper — AI code review poller (ZlayaZanuda bot)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=zanuda
Group=zanuda
WorkingDirectory=/mnt/data/projects/review-helper
# Secrets live in .env, loaded by dotenv at startup.
ExecStart=/usr/local/bin/node dist/index.js
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=NODE_ENV=production
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
