# deploy/zanuda.service.example
[Unit]
Description=Zanuda the Reviewer — AI code review poller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=<SERVICE_USER>
Group=<SERVICE_USER>
WorkingDirectory=<REPO_PATH>
# Point at a local config file with your overrides (allowlist, paths, etc).
# The file is NOT committed — create it from config/default.yaml as a base.
Environment=ZANUDA_CONFIG=<LOCAL_CONFIG_PATH>
Environment=NODE_ENV=production
# Secrets (GITHUB_TOKEN, API keys) live in .env, loaded by dotenv at startup.
ExecStart=<NODE_PATH> dist/index.js
Environment=PATH=<NODE_DIR>:/usr/bin:/bin
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
